# supervisord-local.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=%(here)s/src/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=root
childlogdir=%(here)s/src/logs/


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[group:clients]
programs = server, flower, notebook

[program:server]
command = ./src/manage.py runserver_plus
directory=%(here)s
killasgroup = true
stopasgroup = true
autostart = true
autorestart = true
stdout_logfile = %(here)s/src/logs/server.log
stderr_logfile = %(here)s/src/logs/server.log
priority = 1

[program:flower]
command = celery -A db flower --port=8001 --loglevel=INFO
directory=%(here)s/src
numprocs= 1
stdout_logfile=%(here)s/src/logs/flower.log
stderr_logfile=%(here)s/src/logs/flower.log
autostart = true
autorestart = true
startsecs = 20
stopwaitsecs = 600
killasgroup = true
priority = 999

[program:notebook]
command = ./src/manage.py shell_plus --notebook
directory=%(here)s
numprocs= 1
stdout_logfile=%(here)s/src/logs/notebook.log
stderr_logfile=%(here)s/src/logs/notebook.log
autostart = true
autorestart = true
startsecs = 20
stopwaitsecs = 600
killasgroup = true
priority = 5
